# netbsd/x86/cf/site/disk/backup/root/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.1500 installed on Thu Mar  9 04:32:15 2006)
# (Cron version -- $NetBSD: crontab.c,v 1.26 2005/03/16 02:53:55 xtraeme Exp $)
#
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin
HOME=/var/log
CRON_WITHIN=7200
#
#minute hour    mday    month   wday    command
#
# rotate log files every hour, if necessary
0       *       *       *       *       /usr/bin/newsyslog

